# deployment/telepanorama/decorator/cron/crontab

* * * * * root /usr/bin/wget -O /var/log/decorator.log decorator-nginx

